# src/hil_runtime/pi/systemd/duckpark-dut-result-receiver.service
[Unit]
Description=DuckPark Pi DUT Result Receiver
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/duckpark/src
EnvironmentFile=/etc/duckpark/pi-gateway.env
ExecStart=/bin/bash /opt/duckpark/src/hil_runtime/pi/scripts/start_pi_dut_result_receiver.sh
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
